# extras/ani-notify.service
[Unit]
Description=Aviso de episodios nuevos (ani-cli-mx)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/ani-notify
