# extras/ani-notify.timer
[Unit]
Description=Chequeo periodico de episodios nuevos (ani-cli-mx)

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
